# CondCore/Utilities/o2o/deployed/offline/runTimeBasedO2O.crontab
0 */2 * * * /bin/sh /data/O2O/scripts/runTimeBasedO2O.sh >/dev/null 2>&1 
